# camilladsp.service
[Unit]
Description=CamillaDSP Audio Processor (bit-perfect USB→DAC)
After=sound.target uac2-gadget.service
Wants=sound.target uac2-gadget.service
Requires=uac2-gadget.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/camilladsp-wrapper
Restart=always
RestartSec=3
Nice=-10
LimitRTPRIO=99
LimitMEMLOCK=65536

[Install]
WantedBy=multi-user.target
